# startup/kde-session-units/kdeinit.service.in
[Unit]
Description=KDE Init
Requires=xorg.target dbus.socket lnusertemp@cache.service lnusertemp@socket.service lnusertemp@tmp.service
After=xorg.target dbus.socket lnusertemp@cache.service lnusertemp@socket.service lnusertemp@tmp.service ksplashx.service


[Service]
Type=dbus
BusName=org.kde.klauncher
ExecStart=@KDE4_LIBEXEC_INSTALL_DIR@/start_kdeinit_wrapper +kcminit_startup
# Do all of these need to be added to every user service unit file?
# Need to set DISPLAY as long as we're running on :1 alongside a production session
Environment=DISPLAY=:1 LD_BIND_NOW=true
#Environment=XDG_RUNTIME_DIR=/run/user/%I
#Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%i/dbus/user_bus_socket
#The wrapper forks and exits,  so RemainAfterExit is needed
RemainAfterExit=true

[Install]
WantedBy=systemdkde.target

