# startup/kde-session-units/ksplashx.service.in
[Unit]
Description=Desktop Splash Screen
After=xorg.target
Requires=xorg.target

[Service]
ExecStart=@KDE4_BIN_INSTALL_DIR@/ksplashx Default
Environment=DISPLAY=:1
