# startup/kde-session-units/lnusertemp@.service.in
[Unit]
Description=KDE session symlinks

[Service]
ExecStart=@KDE4_LIBEXEC_INSTALL_DIR@/lnusertemp %I
Type=oneshot
