# startup/kde-session-units/systemdkde.target
[Unit]
Description=The KDE Desktop Environment (native systemd session)
Requires=kdeinit.service
AllowIsolate=true
