# startup/kde-session-units/systemdkde.target
[Unit]
Description=The KDE Desktop Environment (native systemd session)
Wants=ksmserver.service ksplashx.service
AllowIsolate=true

